# files/etc/init/redash_web.conf
description "re:dash web"
chdir /opt/redash/redash/current

kill signal INT

respawn
respawn limit 10 5

start on started network-services
stop on stopping network-services

script
  exec 2>>/var/log/redash.log
  set -x
  . /opt/redash/bin/activate
  . /opt/redash/.env
  ./manage.py runserver
end script
